# src/stats/supervisord.conf.template
[supervisord]
nodaemon=true
user=root

[program:telegraf]
command=/usr/bin/telegraf
